# installer/systemd/zigbee2mqtt.service
[Unit]
Description=Zigbee2MQTT Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=smarthub
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start --prefix /opt/smart-hub/zigbee2mqtt
WorkingDirectory=/opt/smart-hub/zigbee2mqtt
Restart=always
RestartSec=5
StandardOutput=append:/var/log/smart-hub/zigbee2mqtt.log
StandardError=append:/var/log/smart-hub/zigbee2mqtt.error.log

[Install]
WantedBy=multi-user.target